# alive_service.service.j2
[Unit]
Description=Alive Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 {{ script_path }}
Restart=always
User={{ user }}

[Install]
WantedBy=multi-user.target